# twitchtest/supervisord.conf
[supervisord]
nodaemon=true

[program:downloadtwitchnumbers]
command = sh -c "sleep 59;python3 /usr/bin/getviewersandsavetodb.py"
autorestart = true
startsecs = 0
exitcodes = 0